# conf/service/fcgiwrap/finish
#!/usr/bin/execlineb -P
rm -f /run/fcgiwrap/fcgiwrap.sock
